# private_dot_config/systemd/user/sunshine.service.d/override.conf
[Service]
ExecStartPre=
ExecStartPre=%h/bin/wake-display

[Install]
WantedBy=default.target
