# ignore-etc/glamdring/systemd/system/updater.service
[Unit]
Description="A script to update system"

[Service]
ExecStart=/opt/updater

[Install]
WantedBy=multi-user.target
