# private_dot_config/systemd/user/esde-backup.service
[Unit]
Description="A script to backup ES-DE configs and data."

[Service]
Type=oneshot
ExecStart=/home/gandalf/bin/ESDE-backup

