# private_dot_config/systemd/user/optimize-ra-screenshots.service
[Unit]
Description=Optimize RetroArch screenshots

[Service]
Type=oneshot
ExecStart=%h/bin/optimize-ra --screenshots-only
WorkingDirectory=%h
