# private_dot_config/systemd/user/optimize-ra-input.service
[Unit]
Description=Clean stale joypad bindings from retroarch.cfg

[Service]
Type=oneshot
ExecStart=%h/bin/optimize-ra --input-only
