# private_dot_config/systemd/user/optimize-ra-input.path
[Unit]
Description=Watch retroarch.cfg for changes

[Path]
PathModified=%h/.config/retroarch/retroarch.cfg
Unit=optimize-ra-input.service

[Install]
WantedBy=default.target
