# private_dot_config/systemd/user/playerctld.service
[Unit]
Description=playerctld MPRIS daemon

[Service]
ExecStart=/usr/bin/playerctld daemon
Restart=on-failure

[Install]
WantedBy=default.target
