# private_dot_config/systemd/user/screenshot-optimize.service
[Unit]
Description=Optimize RetroArch Screenshots

[Service]
Type=oneshot
ExecStart=%h/bin/optimize-ra_screenshots
WorkingDirectory=%h
