# ignore-etc/proxmox/systemd/container-recovery.timer
[Unit]
Description="5 Min Timer for Container Recovery Service"

[Timer]
OnCalendar=*:0/5
AccuracySec=1s

[Install]
WantedBy=timers.target
